# supervisord.conf
[supervisord]
nodaemon=true

[program:softioc]
command=/usr/bin/procServ -q --name=softioc -i ^C^D^] -c "softIoc -d /tmp/test.db" 2048 
